# root/etc/systemd/system/fail2ban.service.d/nethserver.conf
[Unit]
After=network.target iptables.service shorewall.service

[Service]
ExecStartPre=/sbin/e-smith/expand-template /etc/fail2ban/jail.local

