# root/etc/systemd/system/fail2ban.service.d/nethserver.conf
[Unit]
After=network.target iptables.service shorewall.service


