# root/etc/systemd/system/shorewall.service.d/nethserver.conf
[Unit]
# reload fail2ban when shorewall is reloaded
PropagatesReloadTo=fail2ban.service
