# root/etc/cron.d/fail2ban-statistics
# save all the bans of the day to the json statistics

45 23 * * * root /usr/libexec/nethserver/fail2ban-statistics
